# crontab
* 3 * * * /home/recovery/qnx_backup_and_recovery/backup.sh 

