# crontab
* 1 * * * /opt/qnx_backup_and_recovery/backup.sh 

